# cnf/supervisor.ini
[supervisord]
user=nobody
pidfile=/var/www/vircuum/tmp/logs/supervisord.pid
logfile=/var/www/vircuum/tmp/logs/supervisor.log

[supervisorctl]


[program:cron]
command=/var/www/vircuum/env/bin/python -u /var/www/vircuum/tasks.py cron
stdout_logfile=/var/www/vircuum/tmp/logs/cron.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment = ENVIRONMENT=development,APPLICATION=vircuum

